# docs/examples/host1/home/app5-vm/etc/mesh-init/pmond.socket
[Socket]
ListenStream=/home/app5/run/pmond/control.sock
Accept=false
SocketMode=0660
