# crates/mesh-init/testdata/activated_svc.service
[Service]
ExecStart = "/bin/sh -c 'echo SUCCESS'"
OOMScoreAdjust = -800
