# crates/mesh-init/testdata/activated_svc.socket
[Socket]
ListenStream=14032
Accept=true
