# crates/mesh-init/testdata/auth_test_svc.socket
[Socket]
ListenStream=/tmp/mesh-init-auth-test.sock
ListenStream=14033
Accept=true
