# crates/mesh-init/testdata/auth_test_svc.service
[Service]

ExecStart = "/bin/sh -c 'echo SUCCESS PEER=$X_PEER_UID'"
OOMScoreAdjust = -800
[[Peer]]
uid = 9999
