# crates/mesh-init/testdata/mesh9p-listener.socket
[Socket]
ListenStream=14010
Accept=false
