# crates/mesh-init/testdata/mesh9p-listener-uds.socket
[Socket]
ListenStream=/tmp/mesh9p_listener.sock
Accept=false
