# crates/mesh-init/testdata/mesh9p.socket
[Socket]
ListenStream=14019
Accept=true
