# crates/mesh-init/testdata/pmond.socket
[Socket]
ListenStream=14090
Accept=false
