# crates/mesh-init/testdata/ssh-mesh.socket
[Socket]
ListenStream=14022
Accept=false
