# docs/examples/host2/home/system/etc/mesh-init/mesh9p-host2.socket
[Socket]
Service=mesh9p
ListenStream=15102
Accept=true
