# docs/examples/host2/home/system/etc/mesh-init/pmond.service
# Process monitor for host2.
#
# pmond exposes a JSONL API over UDS at:
#
#   $HOME/run/pmond/control.sock
#
# ssh-mesh exposes HTTP routes that forward to this JSONL socket.

[Service]
# Socket activation is declared in pmond.socket.

ExecStart = "pmond"
OOMScoreAdjust = -820
[Resources]
memory_high = "128M"
memory_max = "256M"
cpu_weight = 50

[Environment]
PATH = "/opt/ssh-mesh/bin:/opt/busybox/bin"
RUST_LOG = "info"
