# docs/examples/host2/home/system/etc/mesh-init/ssh-mesh.service
# mesh-init service definition for host2's ssh-mesh node.
#
# The start script copies this file to:
#
#   $HOME/etc/mesh-init/ssh-mesh.toml
#
# That is mesh-init's normal host1-layout config path. ssh-mesh also uses its
# normal defaults:
#
#   $HOME/.ssh                 key material
#   $HOME/etc/ssh-mesh     mesh.yaml and per-host1 authorized_keys
#   $HOME/run/ssh-mesh        mux sockets and local runtime sockets

[Service]

ExecStart = "ssh-mesh"
OOMScoreAdjust = -900
[[Peer]]
uid = 103

[Resources]
# These are advisory in unprivileged bubblewrap examples. They are included so
# the config exercises mesh-init's resource model when cgroup access is present.
memory_high = "256M"
memory_max = "512M"
cpu_weight = 100

[Environment]
# Include the expected install location first, while still allowing a host1's
# normal PATH to provide the binaries.
PATH = "/opt/ssh-mesh/bin:/opt/busybox/bin"
RUST_LOG = "info"

# Fixed and unique TCP listeners for the example mesh. These are env vars
# because ssh-mesh currently reads listener ports from SSH_PORT/HTTP_PORT during
# startup even when mesh.yaml also contains port fields.
SSH_PORT = "18222"
HTTP_PORT = "18280"

# Enable OpenSSH ControlMaster-compatible mux sockets for maintained client
# connections. sshmc can use these sockets from inside the node.
MUX_DIR = "/home/system/run/ssh-mesh/mux"

# Host2 runs as uid 0 inside the bwrap host1 namespace. ssh-mesh's root default
# for mesh-init is /run/mesh-init/control.sock, while this example intentionally
# keeps mesh-init in the HOME-relative host1 layout, so this one override is
# needed.
MESH_INIT_SOCK = "/home/system/run/mesh-init/control.sock"

# Terminal/host1 examples check /home/USER. Inside the sandbox /home is private
# example state, not the host /home.
SSH_MESH_HOME_ROOT = "/home"
