# docs/examples/host1/home/app1-bwrap/etc/mesh-init/pmond.socket
[Socket]
ListenStream=/home/app1/run/pmond/control.sock
Accept=false
SocketMode=0660
