# docs/examples/host3-vm/home/system/etc/mesh-init/activate-app4-ch.socket
[Socket]
ListenStream=/tmp/mesh/shared/app4-ch/trusted.sock
Accept=true
