# docs/examples/host3-vm/home/system/etc/mesh-init/activate-app3-crosvm.socket
[Socket]
ListenStream=/tmp/mesh/shared/app3-crosvm/trusted.sock
Accept=true
