# docs/examples/host3-vm/home/system/etc/mesh-init/activate-app2-qemu.socket
[Socket]
ListenStream=/tmp/mesh/shared/app2-qemu/trusted.sock
Accept=true
