# docs/examples/host3-vm/home/system/etc/mesh-init/mesh9p-host3-vm.socket
[Socket]
Service=mesh9p
ListenStream=15103
Accept=true
