# docs/examples/ssh-mesh-activation-order/ssh-mesh.socket
[Socket]
# Ordered activation without FileDescriptorName.
#
# TCP listeners are claimed by ssh-mesh in this TCP order:
#   1. SSH TCP
#   2. mesh H2/mTLS TCP
#   3. public static HTTP from <ssh-mesh_home>/web
#   4. public static HTTPS from <ssh-mesh_home>/web
#   5. admin/API HTTP
ListenStream=0.0.0.0:15022
ListenStream=0.0.0.0:15443
ListenStream=0.0.0.0:80
ListenStream=0.0.0.0:443
ListenStream=127.0.0.1:15080

# Common mesh JSONL/JSON-RPC control socket:
ListenStream=/home/ssh-mesh/run/ssh-mesh/jsonl.sock

# Trusted SSH transports are claimed after JSONL:
#   UDS    - trusted SSH over UDS
#   VSOCK  - trusted SSH over AF_VSOCK
ListenStream=/home/ssh-mesh/run/ssh-mesh/trusted-ssh.sock
ListenStream=vsock:2:15022

Accept=false
Service=ssh-mesh
