# docs/examples/host1/home/system/etc/mesh-init/activate-app1-bwrap.socket
[Socket]
ListenStream=/tmp/mesh/shared/app1-bwrap/trusted.sock
Accept=true
