# docs/examples/host1/home/system/etc/mesh-init/activate-app5-vm.socket
[Socket]
ListenStream=/tmp/mesh/shared/app5-vm/trusted.sock
Accept=true
