# docs/examples/host1/home/system/etc/mesh-init/mesh9p.socket
[Socket]
ListenStream=15101
Accept=true
