# docs/examples/host1/home/system/etc/mesh-init/pmond.socket
[Socket]
ListenStream=/home/system/run/pmond/control.sock
Accept=false
SocketMode=0660
